# software/cfg/supervisor.ini
[program:news]
command=/home/spot/software/news/news.py
startsec=10
redirect_stderr=true
autorestart=true
autostart=true
stopsignal=INT

[program:weather]
command=/home/spot/software/weather/weather.py
startsec=10
redirect_stderr=true
autorestart=true
autostart=true
stopsignal=INT

[program:clock]
command=/home/spot/software/clock/clock.py
startsec=10
redirect_stderr=true
autorestart=true
autostart=true
stopsignal=INT

[program:spotify]
command=/home/spot/software/spotify/spotify.py
startsec=10
redirect_stderr=true
autorestart=true
autostart=true
stopsignal=INT

[program:webserver]
command=/home/spot/software/web-server/app.py
startsec=10
redirect_stderr=true
autorestart=true
autostart=true
stopsignal=INT